# support/debian/prerm
systemctl stop uploadnewip
systemctl disable uploadnewip
